# app-conf/supervisor/supervisor.conf
[supervisord]
nodaemon=true
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root                    ;
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

[program:php-cycle]
command=/usr/bin/php /var/www/majordomo/cycle.php
stdout_logfile=/var/log/cycle.log
stdout_logfile_maxbytes=1MB
autostart=true
autorestart=true
user=www-data
numprocs=1
stopasgroup=true
stopsignal=QUIT
priority=100

[inet_http_server]
port=0.0.0.0:9001
username=majordomo
password=majordomo

[include]
files = /etc/supervisor/conf.d/*.conf
